# apps/web/supervisord.jail-web.conf
; ============================================================================
; SIBERMAS Web Jail — Supervisor Configuration (FreeBSD Jails)
; ============================================================================
; Deploy to web jail: /usr/local/etc/supervisord.d/sibermas-web.conf
;
; NOTE: This file is for the web jail ONLY (Next.js standalone).
; The API jail (PHP-FPM + queue workers) has its own supervisor config.
; ============================================================================

[program:sibermas-web]
process_name=%(program_name)s
command=/usr/local/bin/node /usr/local/www/sibermas/web/server.js
directory=/usr/local/www/sibermas/web
autostart=true
autorestart=true
user=www
umask=002
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0",UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/web.log
stopwaitsecs=30
