# apps/web/supervisord.jail-web.conf
; ============================================================================
; SIBERMAS Web Jail — Supervisor Configuration (FreeBSD Jails)
; ============================================================================
; Deploy to web jail: /usr/local/etc/supervisord.d/sibermas-web.conf
;
; SCALING MODE: 4 instances (1 per CPU core) untuk 5000 concurrent users.
; Port: 3000 (instance 00), 3001 (01), 3002 (02), 3003 (03).
; Nginx-proxy jail harus punya upstream web_cluster dengan port ini.
;
; NOTE: PORT env pakai prefix "30" + process_num supaya hasilnya jadi
; 3000, 3001, 3002, 3003 (bukan 00, 01, 02, 03 yang privileged).
;
; Untuk single-instance (dev/low-load), set numprocs=1 dan PORT="3000".
; ============================================================================

[program:sibermas-web]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/node /usr/local/www/sibermas/web/server.js
directory=/usr/local/www/sibermas/web
autostart=true
autorestart=true
user=www
numprocs=4
umask=002
environment=NODE_ENV="production",PORT="30%(process_num)02d",HOSTNAME="0.0.0.0",UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/web.log
stopwaitsecs=30
