# conf/supervisord-web-cluster.conf
; SIBERMAS — Web Jail Supervisor (Next.js Cluster x4)
; Deploy: /usr/local/etc/supervisord.d/sibermas-web.conf
; Target: web jail, 4 instance (1 per CPU core)

[program:sibermas-web]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/node /usr/local/www/sibermas/web/server.js
directory=/usr/local/www/sibermas/web
autostart=true
autorestart=true
user=www
numprocs=4
umask=002
environment=NODE_ENV="production",PORT="%(process_num)02d",HOSTNAME="0.0.0.0",UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/web.log
stopwaitsecs=30
