# apps/api/supervisord.conf
; ============================================================================
; SIBERMAS Supervisor Configuration (FreeBSD — Single Server)
; ============================================================================
; Deploy to: /usr/local/etc/supervisord.d/sibermas.conf
; Root app : /usr/local/www/apache24/data/Sibermas2026
;
; JAIL MODE: Gunakan file terpisah per jail:
;   apps/api/supervisord.jail-api.conf   → API jail (queue workers)
;   apps/web/supervisord.jail-web.conf   → Web jail (Next.js standalone)
;
; Binaries di FreeBSD berada di /usr/local/bin — supervisord tidak mewarisi
; $PATH shell interaktif, jadi WAJIB pakai path absolut atau worker akan gagal
; dengan "can't find command 'php'".
;
; Horizon dihapus karena paket `laravel/horizon` belum terinstal di
; composer.json. Semua queue ditangani oleh `queue:work` standar (Laravel 13
; core). Tambahkan `[program:sibermas-horizon]` lagi setelah `composer require
; laravel/horizon` dijalankan.
; ============================================================================

; ----------------------------------------------------------------------------
; Standard Queue Worker (queues: default, critical, high)
; ----------------------------------------------------------------------------
[program:sibermas-worker-default]
process_name=%(program_name)s_%(process_num)02d
; NOTE: queue:work without a connection arg uses QUEUE_CONNECTION from .env.
; This keeps supervisord in sync with the env file (audit finding C-001).
command=/usr/local/bin/php /usr/local/www/apache24/data/Sibermas2026/apps/api/artisan queue:work --queue=default,critical,high --sleep=3 --tries=3 --max-time=3600
directory=/usr/local/www/apache24/data/Sibermas2026/apps/api
autostart=true
autorestart=true
user=www
numprocs=2
; OPS-FreeBSD 2026-05-12: umask 002 so files the worker writes
; (storage/logs/*, storage/framework/cache, etc) stay group-writable
; for both www (php-fpm) and the deploy user running artisan from SSH.
umask=002
environment=UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/worker-default.log
stopwaitsecs=3600

; ----------------------------------------------------------------------------
; Low Priority Queue Worker
; ----------------------------------------------------------------------------
[program:sibermas-worker-low]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /usr/local/www/apache24/data/Sibermas2026/apps/api/artisan queue:work --queue=low,low-priority --sleep=5 --tries=3 --max-time=3600
directory=/usr/local/www/apache24/data/Sibermas2026/apps/api
autostart=true
autorestart=true
user=www
numprocs=1
umask=002
environment=UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/worker-low.log
stopwaitsecs=3600

; ----------------------------------------------------------------------------
; Long Running Job Worker (email notifications, reports, etc.)
; ----------------------------------------------------------------------------
[program:sibermas-worker-long]
process_name=%(program_name)s
command=/usr/local/bin/php /usr/local/www/apache24/data/Sibermas2026/apps/api/artisan queue:work --queue=long --sleep=10 --tries=1 --max-time=7200
directory=/usr/local/www/apache24/data/Sibermas2026/apps/api
autostart=true
autorestart=true
user=www
numprocs=1
umask=002
environment=UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/worker-long.log
stopwaitsecs=3600

; ----------------------------------------------------------------------------
; Next.js Web Frontend (standalone output)
; ----------------------------------------------------------------------------
; Struktur folder setelah `pnpm --filter web build`:
;   apps/web/.next/standalone/         ← root standalone
;   apps/web/.next/standalone/apps/web/server.js   ← entry point
;   apps/web/.next/standalone/apps/web/.next/static (copy-in manually)
;   apps/web/.next/standalone/apps/web/public      (copy-in manually)
; Lihat docs/DEPLOY_FREEBSD.md langkah ke-8.
[program:sibermas-web]
process_name=%(program_name)s
command=/usr/local/bin/node /usr/local/www/apache24/data/Sibermas2026/apps/web/.next/standalone/apps/web/server.js
directory=/usr/local/www/apache24/data/Sibermas2026/apps/web/.next/standalone/apps/web
autostart=true
autorestart=true
user=www
umask=002
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1",UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/web.log
stopwaitsecs=30

; ----------------------------------------------------------------------------
; Worker Group (restart semua sekaligus: supervisorctl restart workers:*)
;
; NOTE: sibermas-web TIDAK dimasukkan ke grup `workers` — restart sibermas-web
; berarti bouncing Next.js cluster (downtime 5-10 detik), yang tidak perlu
; kalau cuma mau kick queue. Restart web terpisah:
;   supervisorctl restart sibermas-web
; ----------------------------------------------------------------------------
[group:workers]
programs=sibermas-worker-default,sibermas-worker-low,sibermas-worker-long
priority=999
