# apps/api/supervisord.jail-api.conf
; ============================================================================
; SIBERMAS API Jail — Supervisor Configuration (FreeBSD Jails)
; ============================================================================
; Deploy to api jail: /usr/local/etc/supervisord.d/sibermas.conf
;
; SCALING MODE: 10+4+2 workers untuk 5000 concurrent users.
; Untuk single-server/low-load: turunkan numprocs ke 2+1+1.
; ============================================================================

[program:sibermas-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /usr/local/www/sibermas/api/artisan queue:work --queue=default,critical,high --sleep=2 --tries=3 --max-time=3600
directory=/usr/local/www/sibermas/api
autostart=true
autorestart=true
user=www
numprocs=10
umask=002
environment=UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/worker-default.log
stopwaitsecs=3600

[program:sibermas-worker-low]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /usr/local/www/sibermas/api/artisan queue:work --queue=low --sleep=3 --tries=3 --max-time=3600
directory=/usr/local/www/sibermas/api
autostart=true
autorestart=true
user=www
numprocs=4
umask=002
environment=UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/worker-low.log
stopwaitsecs=3600

[program:sibermas-worker-long]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /usr/local/www/sibermas/api/artisan queue:work --queue=long --sleep=5 --tries=1 --max-time=7200
directory=/usr/local/www/sibermas/api
autostart=true
autorestart=true
user=www
numprocs=2
umask=002
environment=UMASK="0002"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/worker-long.log
stopwaitsecs=3600

[group:workers]
programs=sibermas-worker-default,sibermas-worker-low,sibermas-worker-long
priority=999
