# apps/api/supervisord.conf
; SIBERMAS Supervisor Configuration
; Queue Workers for Laravel Horizon and general queue processing

; ================================================================================
; Horizon Dashboard Process
; ================================================================================
[program:sibermas-horizon]
process_name=%(program_name)s
command=php /usr/local/www/sibermas/artisan horizon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www
redirect_stderr=true
stdout_logfile=/var/log/sibermas/horizon.log
stopwaitsecs=3600

; ================================================================================
; Standard Queue Worker
; ================================================================================
[program:sibermas-worker-default]
process_name=%(program_name)s_%(process_num)02d
; NOTE: queue:work without a connection arg uses QUEUE_CONNECTION from .env.
; This keeps supervisord in sync with the env file (audit finding C-001).
command=php /usr/local/www/sibermas/artisan queue:work --queue=default,critical,high --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=www
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/sibermas/worker-default.log
stopwaitsecs=3600

; ================================================================================
; Low Priority Queue Worker
; ================================================================================
[program:sibermas-worker-low]
process_name=%(program_name)s_%(process_num)02d
command=php /usr/local/www/sibermas/artisan queue:work --queue=low --sleep=5 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/sibermas/worker-low.log
stopwaitsecs=3600

; ================================================================================
; Long Running Job Worker (for email notifications, reports, etc.)
; ================================================================================
[program:sibermas-worker-long]
process_name=%(program_name)s
command=php /usr/local/www/sibermas/artisan queue:work --queue=long --sleep=10 --tries=1 --max-time=7200
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/sibermas/worker-long.log
stopwaitsecs=3600

; ================================================================================
; Next.js Web Frontend (standalone output)
; ================================================================================
[program:sibermas-web]
process_name=%(program_name)s
command=node /usr/local/www/sibermas/apps/web/.next/standalone/apps/web/server.js
directory=/usr/local/www/sibermas/apps/web/.next/standalone/apps/web
autostart=true
autorestart=true
user=www
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1"
redirect_stderr=true
stdout_logfile=/var/log/sibermas/web.log
stopwaitsecs=30

; ================================================================================
; Worker Restart Strategy
; ================================================================================
[group:workers]
programs=sibermas-horizon,sibermas-worker-default,sibermas-worker-low,sibermas-worker-long,sibermas-web
priority=999
